# fencerules.service
[Unit]

Description=fencerules

After=default.target

[Service]

ExecStart=/root/bin/fencerules.sh

[Install]

WantedBy=default.target
