# wallrules.service
[Unit]

Description=wallrules

After=default.target

[Service]

ExecStart=/root/bin/wallrules.sh

[Install]

WantedBy=default.target
